# debian/alloyd.service
[Unit]
Description=Alloy URL routing daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/alloyd
Restart=on-failure

[Install]
WantedBy=default.target
